# daemon/systemd/templates/birdnet_plotter.timer
[Unit]
Description=BirdNET-stream plotter timer

[Timer]
OnCalendar=*:00
Unit=birdnet_plotter.service

[Install]
WantedBy=basic.target