# daemon/systemd/templates/birdnet_miner.timer
[Unit]
Description=BirdNET-stream miner Timer

[Timer]
OnCalendar=*-*-* *:00
Unit=birdnet_miner.service

[Install]
WantedBy=basic.target